# supervisord.ini
[supervisord]
nodaemon=true
user=root

[program:dovecot]
command=/usr/sbin/dovecot -F
redirect_stderr=false
stdout_capture_maxbytes=50MB
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/var/log/dovecot.log
stderr_capture_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stderr_logfile=/var/log/dovecot.err

[program:spamassassin]
command=/usr/sbin/spamd --siteconfigpath=/etc/mail/spamassassin --max-children 5 --create-prefs --user-config -v -u dovecot
redirect_stderr=false
stdout_capture_maxbytes=50MB
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_logfile=/var/log/spamd.log
stderr_capture_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stderr_logfile=/var/log/spamd.err

[program:haraka]
command=/usr/bin/haraka -c /srv/haraka
redirect_stderr=false
stdout_capture_maxbytes=200MB
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=10
stdout_logfile=/var/log/haraka.log
stderr_capture_maxbytes=200MB
stderr_logfile_maxbytes=200MB
stderr_logfile_backups=10
stderr_logfile=/var/log/haraka.err