# supervisord.ini
[supervisord]
nodaemon=true
user=root

[program:crond]
command=crond -f -L /var/log/cron.log
stdout_logfile=/var/log/cron.log
stderr_logfile=/var/log/cron-err.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile_backups=0

[program:dovecot]
command=/usr/sbin/dovecot -F
stdout_logfile=/var/log/dovecot.log
stderr_logfile=/var/log/dovecot-err.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile_backups=0

[program:spamassassin]
command=/usr/sbin/spamd -s /var/log/spamd.log --siteconfigpath=/etc/mail/spamassassin --max-children 5 --create-prefs --user-config -v -u dovecot
stdout_logfile=/var/log/spamd.log
stderr_logfile=/var/log/spamd-err.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile_backups=0

[program:haraka]
command=/usr/bin/haraka -c /srv/haraka
stdout_logfile=/var/log/haraka.log
stderr_logfile=/var/log/haraka-err.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile_backups=0